# systemd/gitops-sync.service
# systemd/gitops-sync.service
[Unit]
Description=ATT GitOps Synchronization Pipeline (Containerized)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Injects the GITHUB_PAT securely from the host keystore
EnvironmentFile=/etc/containers/gitops.env

# 1. Ensure the latest lightweight git client is available
ExecStartPre=/usr/bin/podman pull docker.io/alpine/git:latest

# 2. Execute the localized pull using the injected PAT
ExecStart=/usr/bin/podman run --rm -v /opt/att-infrastructure:/repo:z docker.io/alpine/git:latest -C /repo pull https://tminus36:${GITHUB_PAT}@github.com/TMinus36/att-fcos-infrastructure.git main

# 3. Apply any downloaded Quadlet/Configuration updates dynamically
ExecStartPost=/usr/bin/systemctl daemon-reload